# systemd/llama-server.service
[Unit]
Description=Llama Server Service (Qwen3.6-27B Phase 2 MTP)
Documentation=file:///home/hermes/club-3090/.cursor/plans/llama.cpp_phase_2_mtp_c9a5b755.plan.md
After=network.target

[Service]
User=hermes
WorkingDirectory=/home/hermes/llama

# Graph optimization remains off (GGML_CUDA_GRAPH_OPT=0) as default for Qwen
Environment="GGML_CUDA_GRAPH_OPT=0"

# Explicit LD_LIBRARY_PATH prevents dynamic linker from loading conflicting libraries in /usr/local/lib
Environment="LD_LIBRARY_PATH=/home/hermes/llama/llama.cpp/build/bin"

# Phase 2 MTP profile:
#   131072 ctx, MTP speculative decoding (n=2), q4_0 KV, jinja template.
#   Maintains strict loopback-only routing on 127.0.0.1:8080.
ExecStart=/home/hermes/llama/llama.cpp/build/bin/llama-server \
  -m /home/hermes/llama/models/unsloth-mtp-q4km/Qwen3.6-27B-Q4_K_M.gguf \
  -c 131072 \
  -b 4096 \
  -ub 1024 \
  -ngl 99 \
  -fa on \
  --cache-type-k q4_0 \
  --cache-type-v q4_0 \
  -np 1 \
  --spec-type draft-mtp \
  --spec-draft-n-max 2 \
  --jinja \
  --reasoning off \
  --reasoning-format deepseek \
  --temp 0.6 \
  --top-p 0.95 \
  --top-k 20 \
  --repeat-penalty 1.0 \
  --port 8080 \
  --host 127.0.0.1

Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
