# systemd/llama-server.service
[Unit]
Description=Llama Server Service (Qwen3.6-27B Heretic IQ4_NL + DFlash)
Documentation=file:///home/hermes/club-3090/docs/%23%20Quickstart%3A%20Qwen%203.md
After=network.target

[Service]
User=hermes
WorkingDirectory=/home/hermes/llama

Environment="GGML_CUDA_GRAPH_OPT=0"
Environment="LD_LIBRARY_PATH=/home/hermes/llama/beellama.cpp/build/bin"

# Heretic production profile (port 8080):
#   DavidAU Heretic IQ4_NL target + DFlash drafter.
#   128K YaRN ctx (scale 4, orig 32768), temp 1.0, turbo4/turbo3_tcq KV.
ExecStart=/home/hermes/llama/beellama.cpp/build/bin/llama-server \
  -m /home/hermes/llama/models/heretic/Qwen3.6-27B-NEO-CODE-HERE-2T-OT-IQ4_NL.gguf \
  --spec-draft-model /home/hermes/llama/models/precision/dflash-draft-3.6-q4_k_m.gguf \
  --spec-type dflash \
  --spec-dflash-cross-ctx 1024 \
  --spec-draft-ngl all \
  -ngl 99 \
  -fa on \
  --kv-unified \
  -b 2048 \
  -ub 256 \
  --ctx-size 131072 \
  --rope-scaling yarn \
  --rope-scale 4 \
  --yarn-orig-ctx 32768 \
  --cache-type-k turbo4 \
  --cache-type-v turbo3_tcq \
  --cache-ram 0 \
  -np 1 \
  --jinja \
  --reasoning off \
  --reasoning-format deepseek \
  --temp 1.0 \
  --top-k 20 \
  --min-p 0.0 \
  --repeat-penalty 1.0 \
  --port 8080 \
  --host 127.0.0.1 \
  --metrics \
  --log-timestamps \
  --log-prefix \
  --log-colors off

Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
